# chatbot/remote-camera.service
[Unit]
Description=Remote Camera Server for Chatbot
After=network.target

[Service]
Type=simple
User=${USER}
WorkingDirectory=/home/${USER}/dev/camera
ExecStart=/usr/bin/python remote_camera_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
